# themes/simple-monochrome/.config/systemd/user/wallpaper-daily.timer
[Unit]
Description=Rotate the wallpaper daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
