# themes/simple-monochrome/.config/systemd/user/wallpaper-daily.service
[Unit]
Description=Apply the daily rotating wallpaper

[Service]
Type=oneshot
ExecStart=%h/.config/hypr/wallpaper-daily.sh
